# scripts/com.microagent.imessage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Zero external deps — uses macOS built-ins (bash, /usr/bin/python3, osascript).

Install with:
  cp scripts/com.microagent.imessage.plist ~/Library/LaunchAgents/
  launchctl load ~/Library/LaunchAgents/com.microagent.imessage.plist

The first time the sender runs, macOS will prompt to allow Automation
control of Messages.app — approve it. Grant Full Disk Access to
/bin/bash (or Terminal) if the watcher can't see the outbox.

Adjust the OUTBOX path below if your repo checkout is not ~/microagent.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.microagent.imessage</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/joey/microagent/scripts/imessage-sender.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>IMSG_OUTBOX</key>
        <string>/Users/joey/microagent/imessage-outbox</string>
        <key>PATH</key>
        <string>/opt/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/joey/microagent/imessage-outbox/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joey/microagent/imessage-outbox/launchd.err.log</string>
</dict>
</plist>
